# fs/etc/s6-overlay/s6-rc.d/privoxy/run
#!/command/execlineb -P

# Start daemon
s6-setuidgid privoxy /usr/sbin/privoxy /etc/privoxy/config2