# fs/etc/s6-overlay/s6-rc.d/tor/run
#!/command/execlineb -P

# Start daemon
s6-setuidgid tor /usr/bin/tor